# tools/k-skill-qa-bot/launchd/org.nomadamas.k-skill-qa-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.nomadamas.k-skill-qa-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>$HOME/.local/share/k-skill-qa-bot/bin/run-qa.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>259200</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>__HOME__/.local/share/k-skill-qa-bot</string>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/k-skill-qa-bot/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/k-skill-qa-bot/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>Nice</key>
    <integer>5</integer>

    <key>LowPriorityIO</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
</content>
</invoke>